# scripts/sv/nginx/run
#!/bin/bash

set -x

sv start php || exit 1

source /etc/envvars

exec 2>&1
exec nginx -g "daemon off;"
